# files/systemd/user/hypr-config-setup.service
[Unit]
Description=Setup Hyprland Configuration
After=network.target

[Service]
ExecStart=/usr/bin/setup-hyprland-config.sh
Restart=on-failure

[Install]
WantedBy=default.target
